# AutoLights.service
[Unit]
Description=AutoLights executable starter
DefaultDependencies=no

[Service]
ExecStart=/home/yigit/Desktop/dezyne-ornek/

[Install]